# recipes-snake/qsweep-server-init/files/qsweep-server.service
[Unit]
Description=GNU qsweep-server startup script

[Service]
ExecStart=/usr/bin/qsweep-server

[Install]
WantedBy=multi-user.target
